# ansible/roles/r730xd-openbao/templates/foundation-openbao.service.j2
# Managed by Ansible — r730xd-openbao role
# Owns the OpenBao container lifecycle. The container itself is configured
# to restart=unless-stopped at the Docker level; this unit exists so that
# openbao-auto-unseal.service can depend on it with systemd ordering.

[Unit]
Description=OpenBao server (Docker Compose)
Requires=docker.service
After=docker.service network-online.target zfs-mount.service
Wants=network-online.target
# Data lives on ZFS; refuse to start if the mount is missing.
RequiresMountsFor={{ openbao_data_dir }}

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ openbao_compose_dir }}
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
